# splashy.service
[Unit]
Description=Splashy - Smart Water Gun Turret System
Documentation=https://github.com/yos4i/splashy
After=network.target sound.target
Wants=network-online.target
Requires=pigpiod.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/splashy
Environment=PYTHONPATH=/opt/splashy
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /opt/splashy/core/splashy_system.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/splashy
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Resource limits
LimitNOFILE=65536
MemoryLimit=512M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=splashy

[Install]
WantedBy=multi-user.target